# docker/multi-process/scripts/mysqld.conf
[program:mysqld]
priority=20
directory=/tmp
command=/usr/bin/mysqld_safe
autostart=false
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
